# etc/supervisord.d/jenkins.ini
[program:jenkins]
environment=JENKINS_HOME="/opt/jenkins"
command=/bin/java -jar /usr/share/jenkins-1.532.3.war
priority=1
startsecs=5
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
