# etc/supervisord.d/jenkins.ini
[program:jenkins]
command=
priority=1
startsecs=5
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
